# trackpoint/cron_command.cron
@reboot /home/leo/.setup/trackpoint/set.sh
